# systemd/rex-scheduler.service
[Unit]
Description=Rex Scheduler — Autonomous cron daemon for AgentBox
After=network-online.target rex.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/jack/agentbox
ExecStart=/home/jack/.nvm/versions/node/v22.21.0/bin/node /home/jack/.nvm/versions/node/v22.21.0/bin/tsx src/scheduler.ts
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=AGENT=rex
Environment=PATH=/home/jack/.nvm/versions/node/v22.21.0/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
